# docker/intersections/rse203_Hastings/nojournal/bin/rsu
#!/bin/sh /etc/rc.common
START=80

start () {
	/usr/local/bin/rsu_PRS &
	/usr/local/bin/rsu_mprsolver &
	/usr/local/bin/rsu_msg_transmitter &
}

stop () {
	killall rsu_PRS
	killall rsu_msg_transmitter
	killall rsu_mprsolver
}

